# roles/pangu-nimbus/files/nimbus.ini
[program:nimbus]
command=/opt/apache-storm-0.9.4/bin/storm nimbus
autostart=false
autorestart=true
startsecs=3
environment=PATH=/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/root/bin
environment=PATH=/usr/lib/jvm/java/bin:%(ENV_PATH)s
environment=PATH=/usr/lib/jvm/java/jre/bin:%(ENV_PATH)s
environment=PATH=/opt/apache-storm-0.9.4/bin:%(ENV_PATH)s
